# sar/files/etc/init.d/racoon
#!/bin/sh /etc/rc.common
# Copyright (C) 2009 Claudio Mignanti

START=80
STOP=80

start() {
	mkdir -p /var/racoon/conf
	cp -r /etc/racoon/* /var/racoon/
	racoon
}

stop() {
	killall racoon
}
