# zzz-ninux-ipkg-x86/files/etc/init.d/custom-ninux-startup
#!/bin/sh /etc/rc.common
START=90
ip addr add 192.168.6.65/28 dev ath0 
ip addr add 192.168.6.81/28 dev ath1 
