# zzz-ninux-ipkg-adm5120/files/etc/init.d/custom-user-startup
#!/bin/sh /etc/rc.common
START=90
#ip addr add 192.168.6.65/28 dev wl0
